# deploy/mgtv-danmaku.service.example
[Unit]
Description=MGTV Danmaku Vote Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mgtv-vote
Group=mgtv-vote
WorkingDirectory=/opt/MangoTV_Danmaku
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/MangoTV_Danmaku/.venv/bin/python server/vote_server.py --config /var/lib/mgtv-danmaku/config.json
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
UMask=0077
ReadWritePaths=/var/lib/mgtv-danmaku

[Install]
WantedBy=multi-user.target
